# clusters/vc/cfg/vclogin/rootfs/etc/profile.d/slurm-env.csh
setenv SQUEUE_FORMAT	"%.18i %.9Q %.9q %.12a %.8u %.8T %.12M %.12l %.20e %.6D %R"
setenv SQUEUE_SORT	"-t,-p,e,S"
setenv SQUEUE_ALL	yes
