# clusters/vx/cfg/vxsched/rootfs/etc/profile.d/site-vars.sh
../../../../../common/provision/env/site-vars.sh