# clusters/experiments.containers/cc/cfg/cc1/rootfs/etc/profile.d/slurm-env.sh
export SQUEUE_FORMAT="%.18i %.9Q %.9q %.12a %.8u %.8T %.12M %.12l %.20e %.6D %R"
export SQUEUE_SORT="-t,-p,e,S"
export SQUEUE_ALL=yes
